# spark-mcp/src/spark_mcp/templates/systemd.service.template
[Unit]
Description=spark-mcp MCP server
After=network.target docker.service

[Service]
Type=simple
User={user}
Group={user}
WorkingDirectory={working_dir}
ExecStart={exec_path}
EnvironmentFile={env_file}
# hf-download.sh and other vllm-docker scripts need user-local tools
# (uvx, uv, docker from rootless install) which live in ~/.local/bin.
# systemd's default PATH does not include that dir.
Environment="PATH={working_dir}/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
